# launchd/com.ssm.investment-analyst.push-pages.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ssm.investment-analyst.push-pages</string>

  <key>WorkingDirectory</key>
  <string>/Users/ssm/Documents/Investment Analyst</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/npm</string>
    <string>run</string>
    <string>push:pages</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>600</integer>

  <key>StandardOutPath</key>
  <string>/Users/ssm/Documents/Investment Analyst/logs/push-pages.launchd.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/ssm/Documents/Investment Analyst/logs/push-pages.launchd.err.log</string>
</dict>
</plist>
